# etc/systemd/system/cronie.service.d/msmtp.conf
[Service]
ExecStart=
ExecStart=/usr/sbin/crond -m '/usr/bin/msmtp -t'
